# udev/rules.d/99-autorun-add.rules
SUBSYSTEM=="block", KERNEL=="sd*", ACTION=="add", RUN+="/home/sqz/opt/udev-autorun/udev-autorun run &"
