# udev/rules.d/99-autorun-remove.rules
SUBSYSTEM=="block", KERNEL=="sd*", ACTION=="remove", RUN+="/home/sqz/opt/udev-autorun/udev-autorun remove &"
